# server/agent-bridge.plist.example
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- Copy to ~/Library/LaunchAgents/agent-bridge.plist and replace the paths. -->
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>agent-bridge</string>
  <key>ProgramArguments</key>
  <array>
    <string>/ABSOLUTE/PATH/TO/bun</string>
    <string>run</string>
    <string>/ABSOLUTE/PATH/TO/server/bridge.ts</string>
  </array>
  <key>WorkingDirectory</key>
  <string>/ABSOLUTE/PATH/TO/server</string>
  <key>EnvironmentVariables</key>
  <dict>
    <key>PATH</key>
    <string>/usr/local/bin:/opt/homebrew/bin:/usr/bin:/bin</string>
  </dict>
  <key>RunAtLoad</key>
  <true/>
  <key>KeepAlive</key>
  <true/>
  <key>StandardOutPath</key>
  <string>/tmp/agent-bridge.log</string>
  <key>StandardErrorPath</key>
  <string>/tmp/agent-bridge.log</string>
</dict>
</plist>
